# deploy/systemd/ensure-backend.user.service.template
[Unit]
Description=ENSURE Flask Gunicorn backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__FLASK_DIR__
Environment=PATH=__PYTHON_BIN_DIR__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONUNBUFFERED=1
ExecStart=__START_SCRIPT__ --no-install --server gunicorn --host 0.0.0.0 --port 8010
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=default.target
